# crontab
*/5 * * * * /var/lib/prometheus/logcollect/rft.sh >/tmp/rft.sh.out
